# android9/overlay/etc/profile.d/halium-9.sh
# Ignore unkown color formats for camera and media
export GST_AMC_IGNORE_UNKNOWN_COLOR_FORMATS=yes
